# deploy/systemd/factory-orchestrator.service
[Unit]
Description=Factory worker (orchestrator)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/factory-vm
EnvironmentFile=/opt/factory-vm/deploy/env
ExecStart=/opt/factory-vm/.venv/bin/python -m services.workers --role orchestrator
Restart=always

[Install]
WantedBy=multi-user.target
