# deploy/systemd/factory-backup.service
[Unit]
Description=Factory backup snapshot runner
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/factory-vm
EnvironmentFile=/opt/factory-vm/deploy/env
ExecStart=/opt/factory-vm/.venv/bin/python -m scripts.ops_backup_schedule run
