# deploy/systemd/factory-retention.timer
[Unit]
Description=Factory retention runner schedule

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300
Unit=factory-retention.service

[Install]
WantedBy=timers.target
