# deploy/systemd/factory-cleanup.service
[Unit]
Description=Factory worker (cleanup)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/factory-vm
EnvironmentFile=/opt/factory-vm/deploy/env
ExecStart=/opt/factory-vm/.venv/bin/python -m services.workers --role cleanup
Restart=always

[Install]
WantedBy=multi-user.target
