# deploy/cron/cleanup_tick.cron
0 * * * * cd /opt/factory-vm && /opt/factory-vm/.venv/bin/python -m services.workers --role cleanup --once
